# rustyjack.service
[Unit]
Description=Rustyjack UI Service
After=local-fs.target network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/root/Rustyjack
ExecStart=/usr/local/bin/rustyjack-ui
Restart=on-failure
RestartSec=2
User=root
Environment=RUSTYJACK_ROOT=/root/Rustyjack
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
