# services/rustyjack-ui.service
[Unit]
Description=Rustyjack UI Service (unprivileged)
After=local-fs.target network.target rustyjackd.service
Wants=network.target rustyjackd.service
Requires=rustyjackd.service

[Service]
Type=simple
WorkingDirectory=/var/lib/rustyjack
ExecStart=/usr/local/bin/rustyjack-ui
Environment=RUSTYJACK_DISPLAY_ROTATION=landscape
Environment=RUSTYJACK_ROOT=/var/lib/rustyjack
Restart=on-failure
RestartSec=2
User=rustyjack-ui
Group=rustyjack-ui
SupplementaryGroups=rustyjack gpio spi
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
LockPersonality=true
ReadWritePaths=/var/lib/rustyjack
ReadOnlyPaths=/run/rustyjack

[Install]
WantedBy=multi-user.target
Alias=rustyjack.service
