# services/rustyjackd.service
[Unit]
Description=Rustyjack Daemon (Hardened)
Documentation=https://github.com/yourusername/rustyjack
Requires=rustyjackd.socket
After=local-fs.target network.target rustyjackd.socket
Wants=network.target

[Service]
Type=notify
NotifyAccess=main
WatchdogSec=20s
ExecStart=/usr/local/bin/rustyjackd
Restart=on-failure
RestartSec=5

# Runtime & state dirs
# NOTE: /run/rustyjack is created by tmpfiles.d/rustyjack.conf (not
# RuntimeDirectory) to avoid conflicts with the socket unit that owns
# the socket file under that path.
StateDirectory=rustyjack
StateDirectoryMode=0770
ConfigurationDirectory=rustyjack
ConfigurationDirectoryMode=0770
Group=rustyjack

# Environment
Environment=RUSTYJACKD_SOCKET=/run/rustyjack/rustyjackd.sock
Environment=RUSTYJACK_ROOT=/var/lib/rustyjack
Environment=RUSTYJACKD_SOCKET_GROUP=rustyjack
Environment=RUSTYJACKD_UPDATE_PUBKEY_FILE=/etc/rustyjack/update_pubkey.ed25519
Environment=RUSTYJACK_WIFI_BACKEND=dbus
Environment=RUSTYJACKD_OPS_PROFILE=appliance
Environment=RUSTYJACKD_OPS_WIFI=true
Environment=RUSTYJACKD_OPS_ETH=true
Environment=RUSTYJACKD_OPS_HOTSPOT=true
Environment=RUSTYJACKD_OPS_PORTAL=true
Environment=RUSTYJACKD_OPS_STORAGE=true
Environment=RUSTYJACKD_OPS_POWER=true
Environment=RUSTYJACKD_OPS_UPDATE=true
Environment=RUSTYJACKD_OPS_SYSTEM=true
Environment=RUSTYJACKD_OPS_DEV=false
Environment=RUSTYJACKD_OPS_OFFENSIVE=true
Environment=RUSTYJACKD_OPS_LOOT=false
Environment=RUSTYJACKD_OPS_PROCESS=false

# Privilege and capability management
User=root
# Keep only required capabilities
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_DAC_OVERRIDE
# Ambient capabilities for rootless operation where possible
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN

# Filesystem access restrictions
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/rustyjack /run/rustyjack /tmp/rustyjack /etc/resolv.conf /usr/local/bin
ReadOnlyPaths=/etc/wpa_supplicant
PrivateTmp=true

# Kernel hardening
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible

# Process restrictions
NoNewPrivileges=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true

# Memory protection
MemoryDenyWriteExecute=true

# Network restrictions (allow all needed for WiFi)
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK AF_PACKET

# Syscall filtering
SystemCallFilter=@system-service @mount
SystemCallFilter=~@clock @debug @module @obsolete @raw-io @reboot @swap
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=1024
LimitNPROC=64
TasksMax=64
CPUQuota=80%
MemoryMax=256M

# Working directory
WorkingDirectory=/var/lib/rustyjack

# Security
PrivateDevices=false
# Note: PrivateDevices=true would prevent access to /dev/sd* for mounting
# We need access to block devices for mount operations

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rustyjackd

[Install]
WantedBy=multi-user.target
