# services/rustyjack-ui.service
[Unit]
Description=Rustyjack UI Service (unprivileged)
After=local-fs.target network.target rustyjackd.service
Wants=network.target rustyjackd.service
Requires=rustyjackd.service

[Service]
Type=simple
WorkingDirectory=/var/lib/rustyjack
ExecStart=/usr/local/bin/rustyjack-ui
Environment=RUSTYJACK_DISPLAY_ROTATION=landscape
Environment=RUSTYJACK_ROOT=/var/lib/rustyjack
Restart=on-failure
RestartSec=2

# Run as unprivileged user
User=rustyjack-ui
Group=rustyjack-ui
SupplementaryGroups=rustyjack gpio spi

# Filesystem access restrictions
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/rustyjack
ReadOnlyPaths=/run/rustyjack
PrivateTmp=true

# Kernel hardening
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible

# Process restrictions
NoNewPrivileges=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true

# Memory protection
MemoryDenyWriteExecute=true

# Network restrictions (UI only needs Unix socket to daemon)
RestrictAddressFamilies=AF_UNIX

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete @debug
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=512
LimitNPROC=32
TasksMax=32
MemoryMax=128M

# Device access (GPIO/SPI for display and buttons)
DeviceAllow=/dev/spidev* rw
DeviceAllow=/dev/gpiochip* rw
DeviceAllow=/dev/mem r

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rustyjack-ui

[Install]
WantedBy=multi-user.target
Alias=rustyjack.service
